# deploy/tmi-tf-wh.service
[Unit]
Description=TMI Terraform Webhook Analyzer
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=tmi-tf
Group=tmi-tf
WorkingDirectory=/opt/tmi-tf-wh
ExecStart=/usr/bin/python3 -m uvicorn tmi_tf.server:app --host 127.0.0.1 --port 8080
Restart=always
RestartSec=5
Environment=TMI_OAUTH_IDP=tmi

[Install]
WantedBy=multi-user.target
